# build/sdnagent/root/usr/lib/systemd/system/yunion-sdnagent-eipgw.service
[Unit]
Description=OneCloud Eipgw service
Documentation=https://docs.yunion.cn
After=network.target

[Service]
Type=simple
User=root
Group=root
PIDFile=/var/run/yunion-sdnagent-eipgw-keepalived.pid
ExecStart=/usr/sbin/keepalived -f /etc/keepalived/eipgw.conf --vrrp --dont-fork --log-detail --pid=/var/run/yunion-sdnagent-eipgw-keepalived.pid --vrrp_pid=/var/run/yunion-sdnagent-eipgw-vrrp.pid
ExecReload=/bin/kill -HUP $MAINPID
WorkingDirectory=/
KillMode=control-group
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
